# 19_systemd-configurator/examples/03_timers/cron-daily.service
[Unit]
Description=Run daily task

[Service]
Type=oneshot
ExecStart=/usr/local/bin/daily-task

[Install]
WantedBy=multi-user.target
