# 19_systemd-configurator/examples/03_timers/cron-daily.timer
[Unit]
Description=Daily task timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=cron-daily.service

[Install]
WantedBy=timers.target
