# 19_systemd-configurator/examples/02_socket/tcpecho@.service
[Unit]
Description=Per-connection Echo (%i)

[Service]
ExecStart=/usr/bin/socat - TCP4:127.0.0.1:7
StandardInput=socket
