# 19_systemd-configurator/examples/05_mounts/data.automount
[Unit]
Description=Automount for /mnt/data

[Automount]
Where=/mnt/data
