# 19_systemd-configurator/examples/slice/resource.slice
[Unit]
Description=Resource control slice

[Slice]
CPUQuota=50%
MemoryMax=512M
IOReadBandwidthMax=/dev/sda 50M
