# 19_systemd-configurator/examples/07_timer/hello-every10s.service
[Unit]
Description=Log hello world with timestamp (service for timer)

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo "$(date -Is) hello world"'
