# 19_systemd-configurator/examples/04_paths/path-changed.path
[Unit]
Description=Trigger on file changes

[Path]
PathChanged=/etc/hosts
Unit=reloader.service

[Install]
WantedBy=multi-user.target
