# 19_systemd-configurator/examples/socket/echo.socket
[Unit]
Description=Echo Socket (TCP 9999)

[Socket]
ListenStream=9999
Accept=false
Service=echo.service

[Install]
WantedBy=sockets.target
