# 19_systemd-configurator/examples/09_path/conditional.path
[Unit]
Description=Start conditional.service when /etc/myapp/feature.flag is created

[Path]
PathExists=/etc/myapp/feature.flag
Unit=conditional.service

[Install]
WantedBy=multi-user.target
