# 19_systemd-configurator/examples/service/advanced.service
[Unit]
Description=Advanced options demo
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/hosts

[Service]
Type=simple
ExecStart=/usr/bin/env bash -c 'echo start; sleep 10'
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=10s
User=www-data
Group=www-data
WorkingDirectory=/var/run
Environment=ENV=prod LOG_LEVEL=info
EnvironmentFile=-/etc/default/advanced
LimitNOFILE=65536
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SYS_TIME
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/run

[Install]
WantedBy=multi-user.target
