# 19_systemd-configurator/examples/mnt-example.mount
[Unit]
Description=Mount /mnt/example

[Mount]
What=/dev/sdb1
Where=/mnt/example
Type=ext4
Options=defaults,noatime
TimeoutSec=30

[Install]
WantedBy=multi-user.target
