# 19_systemd-configurator/examples/notify.service
[Unit]
Description=Notify type service demo

[Service]
Type=notify
NotifyAccess=all
ExecStart=/usr/lib/demo-notifier
