# 19_systemd-configurator/examples/oneshot.service
[Unit]
Description=Oneshot service demo

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/touch /tmp/oneshot.done
