# 19_systemd-configurator/examples/reloader.service
[Unit]
Description=Reload something on change

[Service]
Type=oneshot
ExecStart=/bin/systemctl reload some.service
