# 19_systemd-configurator/examples/touch.path
[Unit]
Description=Watch for /tmp/trigger file

[Path]
PathExists=/tmp/trigger
Unit=touch.service

[Install]
WantedBy=multi-user.target
