# 19_systemd-configurator/examples/touch.service
[Unit]
Description=Touch handler

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'date >> /tmp/trigger.log'
